# home/dot_local/etc/clamav/clamdscan.plist.tmpl
{{- if eq .host.distro.family "darwin" -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>clamav.clamdscan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/clamdscan</string>
    <string>/</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>UserName</key>
  <string>root</string>
  <key>StandardErrorPath</key>
  <string>/var/log/clamav/clamdscan.error.log</string>
  <key>StandardOutPath</key>
  <string>/var/log/clamav/clamdscan.output.log</string>
</dict>
</plist>
{{- end -}}
