# home/Library/LaunchAgents/com.juicefs.home.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.juicefs.home</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/juicefs</string>
      <string>mount</string>
      <string>private</string>
      <string>{{ .host.homeParentFolder }}/{{ .user.name }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/juicefs_setup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/juicefs_setup.err</string>
  </dict>
</plist>
