# system/Library/LaunchDaemons/clamav.clamdscan.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>clamav.clamdscan</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ lookPath "clamdscan" }}</string>
    <string>-i</string>
    <string>-r</string>
    <string>--move=/quarantine</string>
    <string>/</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>5</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>UserName</key>
  <string>root</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/clamdscan.error.log</string>
  <key>StandardOutPath</key>
  <string>/Library/Logs/clamdscan.debug.log</string>
</dict>
</plist>
