# home/Library/LaunchDaemons/com.github.muesli.gitomatic.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>EnvironmentVariables</key>
    <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
  </dict>
  <key>Disabled</key>
  <false/>
  <key>Label</key>
  <string>com.github.muesli.gitomatic</string>
  <key>Program</key>
  <string>/usr/local/bin/gitomatic-service</string>
  <key>RunAtLoad</key>
  <true/>
  <key>UserName</key>
  <string>root</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/com.github.muesli.gitomatic.error.log</string>
  <key>StandardOutPath</key>
  <string>/Library/Logs/com.github.muesli.gitomatic.debug.log</string>
</dict>
</plist>
