# home/Library/LaunchDaemons/com.apple.automatedupdates.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
       <key>Label</key>
       <string>com.apple.automatedupdates</string>
       <key>ProgramArguments</key>
       <array>
           <string>/usr/sbin/softwareupdate</string>
           <string>-i -a</string>
       </array>
       <key>StartCalendarInterval</key>
       <dict>
           <key>Minute</key>
           <integer>0</integer>
       </dict>
       <key>StandardErrorPath</key>
	   <string>/var/log/com.apple.automatedupdates.plist.error.log</string>
	   <key>StandardOutPath</key>
	   <string>/var/log/ccom.apple.automatedupdates.plist.debug.log</string>
   </dict>
</plist>