# home/Library/LaunchDaemons/rclone.user.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>rclone.public</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/rclone-mount</string>
		<string>user</string>
		<string>{{ .user.username }}</string>
		<string>s3</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
    <key>StandardErrorPath</key>
	<string>/Users/{{ .user.username }}/.local/share/rclone/error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/{{ .user.username }}/.local/share/rclone/debug.log</string>
</dict>
</plist>